# alpine/lxmd.initd
#!/sbin/openrc-run

description="LXMF Router Daemon"

command="/opt/reticulum/bin/lxmd"
command_args="--service --config /etc/lxmd --rnsconfig /etc/reticulum"
command_user="reticulum:reticulum"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net rnsd
}

start_pre() {
    checkpath --directory --owner root:root --mode 0777 /etc/reticulum/storage
}
