# alpine/rnsd.initd
#!/sbin/openrc-run

description="Reticulum Network Stack Daemon"

command="/opt/reticulum/bin/rnsd"
command_args="--service --config /etc/reticulum"
command_user="reticulum:reticulum"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner reticulum:reticulum --mode 0755 /etc/reticulum/storage
    checkpath --directory --owner reticulum:reticulum --mode 0755 /etc/reticulum/interfaces
}
