# debian/lxmd.service
[Unit]
Description=LXMF Router Daemon
After=network.target rnsd.service
Requires=rnsd.service
Wants=network.target

[Service]
Type=simple
User=reticulum
Group=reticulum
ExecStart=/opt/reticulum/bin/lxmd --config /etc/lxmd --rnsconfig /etc/reticulum
Restart=on-failure
RestartSec=10s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/etc/lxmd /etc/reticulum
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
PrivateDevices=no

[Install]
WantedBy=multi-user.target
