# alpine/lxmd/lxmd.initd
#!/sbin/openrc-run

description="LXMF Router Daemon"

command="/usr/lib/rnsd/bin/lxmd"
command_args="--service --config /etc/lxmd --rnsconfig /etc/reticulum"
command_user="lxmd:lxmd"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
    need rnsd
}

start_pre() {
    checkpath --directory --owner root:root --mode 0777 /etc/reticulum/storage
}